# agent/systemd/osqueryd.service.d/aigis-zero.conf
[Unit]
Before=aigis-zero.service

[Service]
LimitNOFILE=65536
LimitNPROC=4096
LimitMEMLOCK=134217728
Restart=always
RestartSec=3s
TimeoutStartSec=60
