# agent/systemd/osqueryd.service.d/aigis-zero.conf
# osqueryd resource tuning for EDR workloads.
# This drop-in applies resource limits only — it does NOT create any
# ordering or dependency relationship with aigis-zero.service.

[Service]
LimitNOFILE=65536
LimitNPROC=4096
LimitMEMLOCK=134217728
Restart=always
RestartSec=3s
TimeoutStartSec=60
