# agent/osquery/osquery.conf
{
  "options": {
    "host_identifier": "uuid",
    "schedule_splay_percent": 10,
    "events_expiry": 3600,
    "events_max": 500000,
    "disable_audit": "true",
    "enable_bpf_events": "true",
    "disable_events": "false"
  },

  "decorators": {
    "load": [
      "SELECT uuid AS host_uuid FROM system_info;",
      "SELECT hostname FROM system_info;",
      "SELECT version AS osquery_version FROM osquery_info;",
      "SELECT version AS os_version, platform AS os_platform FROM os_version;"
    ],
    "always": [
      "SELECT user AS active_user FROM logged_in_users WHERE user != '' AND type = 'user' ORDER BY time DESC LIMIT 1;"
    ],
    "interval": {
      "3600": ["SELECT total_seconds AS uptime_seconds FROM uptime;"]
    }
  },

  "schedule": {

    "bpf_process_events": {
      "query": "SELECT tid, pid, parent, uid, gid, cid, exit_code, probe_error, syscall, path, cwd, cmdline, duration, ntime, time FROM bpf_process_events;",
      "interval": 30,
      "removed": false,
      "platform": "linux",
      "denylist": false,
      "description": "Process executions via eBPF (execve/fork/clone). cid=cgroup for containers."
    },

    "bpf_socket_events": {
      "query": "SELECT tid, pid, parent, uid, gid, cid, exit_code, probe_error, syscall, local_address, remote_address, local_port, remote_port, time FROM bpf_socket_events;",
      "interval": 30,
      "removed": false,
      "platform": "linux",
      "denylist": false,
      "description": "Network connections via eBPF (bind/connect/accept)."
    },

    "running_processes": {
      "query": "SELECT pid, ppid, name, path, cmdline, state, uid, gid, euid, egid, start_time, resident_size, user_time, system_time, disk_bytes_read, disk_bytes_written, on_disk FROM processes;",
      "interval": 60,
      "description": "All running processes with resource usage."
    },

    "deleted_binaries": {
      "query": "SELECT pid, name, path, cmdline, uid, gid FROM processes WHERE on_disk = 0;",
      "interval": 60,
      "snapshot": true,
      "removed": false,
      "description": "Processes with deleted executables — fileless malware indicator."
    },

    "listening_ports": {
      "query": "SELECT pid, port, protocol, family, address, processes.name, processes.path, processes.cmdline FROM listening_ports LEFT JOIN processes USING (pid);",
      "interval": 60,
      "description": "Listening ports with owning process. Changes indicate new services or backdoors."
    },

    "open_sockets": {
      "query": "SELECT pid, fd, socket, family, protocol, local_address, remote_address, local_port, remote_port, state FROM process_open_sockets WHERE family != 1;",
      "interval": 60,
      "description": "Open TCP/UDP sockets per process (excludes UNIX domain sockets)."
    },

    "logged_in_users": {
      "query": "SELECT type, user, tty, host, time, pid FROM logged_in_users;",
      "interval": 60,
      "removed": false,
      "description": "Active login sessions via utmp."
    },

    "memory_mapped_executables": {
      "query": "SELECT pid, start, end, permissions, offset, path FROM process_memory_map WHERE permissions LIKE '%x%' AND path != '' AND path NOT LIKE '/usr/%' AND path NOT LIKE '/lib%' AND path NOT LIKE '/lib64%';",
      "interval": 120,
      "description": "Executable memory regions outside standard paths — shellcode / injected libs."
    },

    "process_open_files": {
      "query": "SELECT pid, fd, path FROM process_open_files WHERE path NOT LIKE '/proc/%' AND path NOT LIKE '/dev/%' ORDER BY pid LIMIT 1000;",
      "interval": 120,
      "description": "Open file descriptors per process (noise-filtered)."
    },

    "cpu_utilization": {
      "query": "SELECT pid, name, uid, user_time, system_time, resident_size, total_size, (user_time + system_time) AS total_cpu_time FROM processes ORDER BY total_cpu_time DESC LIMIT 25;",
      "interval": 60,
      "snapshot": true,
      "description": "Top 25 CPU-consuming processes. Essential for performance debugging."
    },

    "memory_utilization": {
      "query": "SELECT pid, name, uid, resident_size, total_size FROM processes ORDER BY resident_size DESC LIMIT 25;",
      "interval": 60,
      "snapshot": true,
      "description": "Top 25 memory-consuming processes."
    },

    "system_memory": {
      "query": "SELECT memory_total, memory_free, memory_available, swap_total, swap_free, buffers, cached FROM memory_info;",
      "interval": 60,
      "snapshot": true,
      "description": "System-wide memory stats (from /proc/meminfo)."
    },

    "disk_io": {
      "query": "SELECT name, reads_completed, writes_completed, read_bytes, write_bytes, read_time, write_time FROM disk_events;",
      "interval": 60,
      "description": "Disk I/O statistics per device."
    },

    "crontab_snapshot": {
      "query": "SELECT command, path, minute, hour, day_of_month, month, day_of_week FROM crontab;",
      "interval": 300,
      "snapshot": true,
      "description": "Full crontab dump. Changes may indicate persistence."
    },

    "docker_containers": {
      "query": "SELECT id, name, image, image_id, command, created, state, status FROM docker_containers;",
      "interval": 300,
      "description": "Running Docker containers."
    },

    "users_snapshot": {
      "query": "SELECT uid, gid, username, description, directory, shell, type FROM users;",
      "interval": 3600,
      "snapshot": true,
      "description": "All local user accounts."
    },

    "groups_snapshot": {
      "query": "SELECT gid, gid_signed, groupname FROM groups;",
      "interval": 3600,
      "snapshot": true,
      "description": "All local groups."
    },

    "sudoers": {
      "query": "SELECT header, rule_details FROM sudoers;",
      "interval": 3600,
      "description": "Sudoers entries. Changes = potential privilege escalation."
    },

    "authorized_keys": {
      "query": "SELECT users.username, authorized_keys.key, authorized_keys.key_type, authorized_keys.comment, authorized_keys.key_file FROM users JOIN authorized_keys USING (uid);",
      "interval": 3600,
      "description": "SSH authorized keys. New keys = potential backdoor."
    },

    "user_ssh_keys": {
      "query": "SELECT uid, path, encrypted FROM user_ssh_keys;",
      "interval": 3600,
      "description": "Private SSH keys on disk. Unencrypted = critical finding."
    },

    "kernel_modules": {
      "query": "SELECT name, size, used_by, status, address FROM kernel_modules;",
      "interval": 3600,
      "description": "Loaded kernel modules. Unexpected entries = potential rootkit."
    },

    "suid_bins": {
      "query": "SELECT path, username, groupname, permissions, pid_with_namespace FROM suid_bin;",
      "interval": 3600,
      "description": "SUID/SGID binaries. New entries = privilege escalation vector."
    },

    "startup_items": {
      "query": "SELECT name, type, path, args, status FROM startup_items;",
      "interval": 3600,
      "description": "System startup items."
    },

    "systemd_units": {
      "query": "SELECT id, description, load_state, active_state, sub_state FROM systemd_units WHERE active_state = 'active';",
      "interval": 3600,
      "description": "Active systemd units. New entries = persistence indicator."
    },

    "deb_packages_snapshot": {
      "query": "SELECT name, version, source, arch, revision, status FROM deb_packages;",
      "interval": 3600,
      "snapshot": true,
      "platform": "linux",
      "description": "Installed Debian packages."
    },

    "rpm_packages_snapshot": {
      "query": "SELECT name, version, release, arch, epoch, install_time, vendor, package_group, source FROM rpm_packages;",
      "interval": 3600,
      "snapshot": true,
      "platform": "linux",
      "description": "Installed RPM packages."
    },

    "apt_sources": {
      "query": "SELECT name, source, base_uri, release, maintainer, type FROM apt_sources;",
      "interval": 3600,
      "description": "APT repo sources. Malicious repos = supply chain risk."
    },

    "yum_sources": {
      "query": "SELECT name, baseurl, enabled, gpgcheck FROM yum_sources;",
      "interval": 3600,
      "description": "YUM/DNF repository sources."
    },

    "mounts": {
      "query": "SELECT device, device_alias, path, type, flags FROM mounts;",
      "interval": 3600,
      "description": "Mounted filesystems. Unexpected mounts = potential pivot."
    },

    "interface_details": {
      "query": "SELECT interface, mac, ip_mask, flags, type FROM interface_details;",
      "interval": 3600,
      "snapshot": true,
      "description": "Network interface details."
    },

    "routes": {
      "query": "SELECT destination, netmask, gateway, source, flags, interface, mtu FROM routes;",
      "interval": 3600,
      "description": "Kernel routing table."
    },

    "arp_cache": {
      "query": "SELECT address, mac, interface FROM arp_cache;",
      "interval": 300,
      "removed": false,
      "description": "ARP cache. Changes may indicate ARP poisoning."
    },

    "dns_resolvers": {
      "query": "SELECT id, type, address FROM dns_resolvers;",
      "interval": 3600,
      "description": "DNS resolver configuration."
    },

    "fim_critical_files": {
      "query": "SELECT path, sha256 FROM hash WHERE path IN ('/etc/passwd', '/etc/shadow', '/etc/sudoers', '/etc/crontab', '/etc/hosts', '/etc/resolv.conf', '/etc/nsswitch.conf', '/etc/pam.conf', '/etc/ssh/sshd_config', '/etc/fstab', '/etc/ld.so.conf');",
      "interval": 300,
      "description": "SHA-256 hashes of critical config files. Changes = tampering indicator."
    },

    "shadow_hashes": {
      "query": "SELECT username, password_status, last_change, expire FROM shadow;",
      "interval": 3600,
      "snapshot": true,
      "description": "Password hash status. NP = no password = critical finding."
    },

    "system_info_snapshot": {
      "query": "SELECT hostname, cpu_brand, cpu_physical_cores, cpu_logical_cores, physical_memory, hardware_vendor, hardware_model, hardware_serial, computer_name FROM system_info;",
      "interval": 86400,
      "snapshot": true,
      "description": "Hardware inventory."
    },

    "os_version_snapshot": {
      "query": "SELECT name, version, major, minor, patch, platform, platform_like, codename FROM os_version;",
      "interval": 86400,
      "snapshot": true,
      "description": "OS version details."
    },

    "kernel_info": {
      "query": "SELECT version, arguments, path, device FROM kernel_info;",
      "interval": 86400,
      "snapshot": true,
      "description": "Kernel version and boot arguments."
    },

    "last_logins": {
      "query": "SELECT username, tty, pid, type, time, host FROM last;",
      "interval": 86400,
      "snapshot": true,
      "description": "Historical login records from /var/log/wtmp."
    },

    "osquery_health": {
      "query": "SELECT name, interval, executions, output_size, wall_time, (user_time/executions) AS avg_user_ms, average_memory, last_executed FROM osquery_schedule;",
      "interval": 3600,
      "snapshot": true,
      "description": "Performance stats per query. Use to detect expensive queries."
    },

    "uptime": {
      "query": "SELECT total_seconds, days, hours, minutes FROM uptime;",
      "interval": 3600,
      "snapshot": true,
      "description": "System uptime. Unexpected reboots are security events."
    },

    "disk_info": {
      "query": "SELECT disk_index, type, id, disk_size, hardware_model, name, serial, vendor FROM disk_info;",
      "interval": 86400,
      "snapshot": true,
      "description": "Physical disk inventory."
    }
  }
}
